# system/systemd/usb-mount@.service
[Unit]
Description=Mount USB Drive on %i
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=$lbhomedir/sbin/usb-mount.sh add %i
ExecStop=$lbhomedir/sbin/usb-mount.sh remove %i
