# system/cron/cron.hourly/01-setdatetime
#!/bin/sh

# Not used on DietPi
if [ -e /boot/dietpi/.hw_model ]
then
  exit 0
fi

# Update Date/Time every hour
if [ -x $LBHOMEDIR/sbin/setdatetime.pl ]
then
  $LBHOMEDIR/sbin/setdatetime.pl
fi
