# system/cron/cron.hourly/01-setdatetime
#!/bin/sh

# Update Date/Time every hour
if [ -x $LBHOMEDIR/sbin/setdatetime.pl ]
then
  $LBHOMEDIR/sbin/setdatetime.pl
fi
