# system/cron/cron.hourly/01-setdatetime
#!/bin/sh

# Update Date/Time every hour
if [ -x /opt/loxberry/sbin/setdatetime.pl ]
then
  /opt/loxberry/sbin/setdatetime.pl
fi
