# system/cron/cron.hourly/02_log_maint
#!/bin/bash
$LBHOMEDIR/sbin/log_maint.pl action=reduce_notifys
$LBHOMEDIR/sbin/log_maint.pl action=reduce_logfiles
