# system/cron/cron.hourly/03-apachestate
#!/bin/bash
$LBHOMEDIR/sbin/checkapachestate.sh
