# system/cron/cron.hourly/03-healthcheck
#!/bin/bash
$LBHOMEDIR/sbin/healthcheck.pl action=check output=notification
