# system/cron/cron.hourly/db_maint_hourly
#!/bin/bash
$LBHOMEDIR/bin/db_maint.pl action=reduce_notifys
$LBHOMEDIR/bin/db_maint.pl action=reduce_logfiles



