# system/cron/cron.d/setloxberryid
MAILTO=""
PATH=/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

# m h  dom mon dow   command

# Additional run at 23:59 to rotate the statistics file
10 3 */3 * * 	loxberry /opt/loxberry/sbin/setloxberryid.pl -d > /dev/null 2>&1
@reboot 		loxberry /opt/loxberry/sbin/setloxberryid.pl > /dev/null 2>&1
