# system/cron/cron.daily/01-log_maint
#!/bin/bash
$LBHOMEDIR/sbin/log_maint.pl action=backup_logdb
