# system/cron/cron.daily/02-pluginsupdate
#!/bin/bash
$LBHOMEDIR/sbin/pluginsupdate.pl
