# system/cron/cron.daily/03-healthcheck
#!/bin/bash
sleep $[RANDOM%360]m
$LBHOMEDIR/sbin/healthcheck.pl action=check output=notification
