# system/cron/cron.weekly/db_maint_weekly
#!/bin/bash
$LBHOMEDIR/bin/db_maint.pl action=backup_logdb




